# systemd/supabase-restore.timer
[Unit]
Description=Weekly Timer for Supabase Restore

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
