# systemd/supabase-backup.service
[Unit]
Description=Supabase Daily Backup Service
After=network.target

[Service]
Type=oneshot
# Project directory
WorkingDirectory=%h/repos/supabase-backup
# Run the script and pipe output to log file in the project directory
ExecStart=/bin/bash -c "mkdir -p %h/repos/supabase-backup/logs && %h/repos/supabase-backup/backup.sh 2>&1 | tee -a %h/repos/supabase-backup/logs/backup.log"

[Install]
WantedBy=default.target
