# systemd/supabase-restore.service
[Unit]
Description=Supabase Restore to Test Environment
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/repos/supabase-backup
ExecStart=/usr/bin/env bash -c "%h/repos/supabase-backup/restore.sh --test --yes"

[Install]
WantedBy=default.target
