# systemd/supabase-backup.service
[Unit]
Description=Supabase Daily Backup Service
After=network.target

[Service]
Type=oneshot
# Update this to your project directory
WorkingDirectory=/home/user/repos/supabase-backup
# Run the script and pipe output to both journal (stdout) and log file
ExecStart=/bin/bash -c '/home/user/repos/supabase-backup/backup.sh 2>&1 | tee -a /var/log/supabase-backup.log'

[Install]
WantedBy=multi-user.target
