# .backup-fixes/2026-04-22/supabase-restore.service.template
[Unit]
Description=Supabase Restore to Test Environment
After=network.target

[Service]
Type=oneshot
WorkingDirectory=${INSTALL_DIR}
ExecStart=/usr/bin/env bash -c "${INSTALL_DIR}/restore.sh --test --yes"

[Install]
WantedBy=default.target
